# printer.conf
[program:printer]
command = sudo /usr/bin/node /home/pi/printer/print.js
user = pi
stdout_logfile = /home/pi/printer/supervisor.log
redirect_stderr = true